# profile.d/hist.sh
#!/bin/sh

export HISTSIZE=60000
export HISTFILESIZE=60000